# daily_message.service
[Unit]
Description = Daily message

[Service]
ExecStart = /daily-message/start.sh

[Install]
WantedBy = multi-user.target